# systemd/zsnap.timer
[Unit]
Description=Run zsnap.py weekly

[Timer]
OnCalendar=Wed *-*-* 03:00:00
Persistent=true
Unit=zsnap.service

[Install]
WantedBy=timers.target
